# chef/atc/templates/default/mount-udev.rules.erb
# File generated by chef, do not edit

# Start <%= @service %> on vagrant mount
SUBSYSTEM=="bdi",ACTION=="add",RUN+="/etc/init.d/<%= @service %> start"

# Stop <%= @service %> on vagrant unmount
SUBSYSTEM=="bdi",ACTION=="remove",RUN+="/etc/init.d/<%= @service %> stop"
